# pacman/unitfiles/clean-pacman-cache.timer
[Unit]
Description=Clean pacman cache monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
