# pacman/unitfiles/clean-pacman-cache.service
[Unit]
Description=Clean pacman cache

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Sc
